# target/linux/ramips/mt7621/base-files/etc/hotplug.d/iface/99-set-irq-affinity
#!/bin/sh
case "$ACTION" in
ifup)
sleep 2
/etc/init.d/set-irq-affinity restart
logger -t "hotplug-irq" "Interface $INTERFACE ($DEVICE) is up, applying IRQ affinity."
;;
esac

